# .config/autostart/autostart.sh.desktop
[Desktop Entry]
Type=Application
Exec=/home/poiuto/Documents/autostart/autostart.sh
Name=autostart.sh
