# .config/autostart/autostart.sh.desktop
[Desktop Entry]
Type=Application
Exec=sh -c '$HOME/.config/remapkey/autostart.sh'
Name=autostart.sh
