# .config/autostart/remapkey.sh.desktop
[Desktop Entry]
Type=Application
Exec=sh -c '$HOME/.config/scripts/autostart/remapkey.sh'
Name=remapkey.sh
