# .config/autostart/remapkey.sh.desktop
[Desktop Entry]
Exec=sh -c '$HOME/.config/scripts/autostart/remapkey.sh'
Name=remapkey.sh
Type=Application
